# rpi/systemd/mnt-wmb\x2dbackup.automount
[Unit]
Description=Automount WatchMyBirds USB backup volume
# Lazy mount: triggered on first access to /mnt/wmb-backup.
# Stick can be inserted/removed at any time without systemd errors.
DefaultDependencies=no
ConditionPathExists=/dev/disk/by-label/WMB-BACKUP

[Automount]
Where=/mnt/wmb-backup
# Idle unmount after 60s of no activity -- lets the stick spin down
# and makes "physical removal while idle" safe.
TimeoutIdleSec=60

[Install]
WantedBy=multi-user.target
