# rpi/systemd/mnt-wmb\x2dbackup.mount
[Unit]
Description=WatchMyBirds USB backup volume
# Stick is fully optional. Boot must not block on it.
DefaultDependencies=no
Conflicts=umount.target
Before=umount.target

[Mount]
What=/dev/disk/by-label/WMB-BACKUP
Where=/mnt/wmb-backup
Type=ext4
# Allowlist: only ext4 is acceptable. Anything else fails the mount,
# UI surfaces the wrong-fs state. Hardlinks (rsync --link-dest)
# require ext4 -- FAT32/exFAT would silently break dedupe.
#
# X-mount.* is handled by util-linux mount(8): it fixes the mounted
# stick root to the unprivileged app user after mount, which makes a
# fresh ext4 stick usable without a separate manual chown step.
Options=rw,nosuid,nodev,noexec,noatime,X-mount.owner=watchmybirds,X-mount.group=watchmybirds,X-mount.mode=0755
TimeoutSec=10

[Install]
WantedBy=multi-user.target
