# rpi/systemd/wmb-wifi-setup.path
[Unit]
Description=Watch for WatchMyBirds WiFi Pending Config

[Path]
PathExists=/var/lib/watchmybirds/pending_wifi.conf

[Install]
WantedBy=multi-user.target
