# rpi/systemd/wmb-wifi-watchdog.timer
[Unit]
Description=Run WiFi Watchdog every 2 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
