# rpi/systemd/wmb-boot-diag.service
[Unit]
Description=WatchMyBirds Boot Diagnostics
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wmb-boot-diag.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
