# rpi/systemd/wmb-wifi-setup.path
[Unit]
Description=Watch for WatchMyBirds WiFi Pending Config

[Path]
PathExists=/opt/app/data/pending_wifi.conf

[Install]
WantedBy=multi-user.target
