# rpi/systemd/wmb-wifi-watchdog.service
[Unit]
Description=WatchMyBirds WiFi Connection Watchdog
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wmb-wifi-watchdog.sh
